# Scripts/youpi.service
[Unit]
Description=youpi start

[Service]
Type=oneshot
ExecStart=/home/pi/YouPi/Scripts/youpi.sh

[Install]
WantedBy=multi-user.target
