# Scripts/youpi.service
[Unit]
Description=youpi start
User=pi

[Service]
Type=oneshot
ExecStart=/home/pi/Youpi/webflask/youpi.py

[Install]
WantedBy=multi-user.target
